# systemd-timer/godaddy-ddns.timer
[Unit]
Description=Set DNS record for GoDaddy

[Timer]
OnBootSec=1m
OnUnitActiveSec=15m

[Install]
WantedBy=timers.target
